# Configs/systemd/user/mango.service
[Unit]
Description=MangoWM Wayland Compositor

[Service]
ExecStart=/usr/bin/mango
Restart=on-failure

[Install]
WantedBy=default.target
